# pi/crontab
*/5 * * * * /home/pi/login.sh
*/1 * * * * sudo /home/pi/tempControl.sh